# services/wifi-antenna-check.service
[Unit]
Description=WiFi antenna health check
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wifi-antenna-check
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wifi-antenna-check
